# startup_runtime/brushless.timer
[Unit]
Description=Run brushless motor sweep every 30 seconds

[Timer]
OnBootSec=60
OnUnitActiveSec=30s

[Install]
WantedBy=timers.target